# bbb/buildroot/board/vaporsound/rootfs-overlay/etc/service/mpd/run
#!/bin/sh

exec /usr/bin/mpd --no-daemon
