# bbb/buildroot/board/vaporsound/rootfs-overlay/etc/service/shairport/run
#!/bin/sh

exec /usr/bin/shairport-sync
