# bbb/buildroot/board/vaporsound/rootfs-overlay/etc/service/pulseaudio/run
#!/bin/sh

exec /usr/bin/pulseaudio --system
